# build/fcm-service.service
[Service]
User=root
Group=root
WorkingDirectory=/root/go/src/fcm-service
ExecStart=/root/go/src/fcm-service/app.exe
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

Description=fcm-service
After=network.target